# slack-canvas-creator-from-threads.service
[Unit]
Description=Slack Canvas Creator from Threads
Documentation=https://github.com/jjj999/slack-canvas-creator-from-threads
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=slackapp
Group=slackapp
WorkingDirectory=/opt/slack-canvas-creator-from-threads
EnvironmentFile=/opt/slack-canvas-creator-from-threads/.env
ExecStart=/opt/slack-canvas-creator-from-threads/venv/bin/python run.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slack-canvas-creator

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/slack-canvas-creator-from-threads
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
